# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:chef-server]
command=run.sh

[program:berkshelf-api]
command=berks-api -c ~/.berkshelf/api-server/config.json
